# config/systemd/user/elephant.service
[Unit]
Description=Omarchy Elephant (Walker backend)
After=graphical-session.target

[Service]
ExecStart=/usr/bin/elephant service run
Restart=on-failure
RestartSec=3
Slice=session.slice

[Install]
WantedBy=default.target
